# etc/systemd/system/unbound-anchor.service.d/unbound-anchor.conf
[Service]
ExecStart=/usr/sbin/unbound-anchor -a /var/lib/unbound/root.key -c /etc/unbound/icannbundle.pem -f /etc/si-resolv.conf -R
