# etc/cron.d/audit-log-rotator
01  *  *  *  *  root  /usr/libexec/audit-log-rotator
