# etc/cron.d/auto-reboot
45 06  *  *  4  root  /usr/local/sbin/rbc -r
