# etc/cron.d/dnf
30 05  *  *  *  root  /usr/local/libexec/dnf.upgrade
