# etc/cron.d/dnf
30 01  *  *  *  root  /usr/local/libexec/dnf.upgrade
