# etc/systemd/system/php81-php-fpm.service.d/php-fpm.conf
[Unit]
RequiresMountsFor=/data
RequiresMountsFor=/nmnh
Before=nginx.service
After=network.target network-online.target ipsec.service remote-fs.target nss-lookup.target

[Service]
TimeoutStopSec=120
