# etc/systemd/system/nginx.service.d/nginx-override.conf
[Unit]
After=network.target network-online.target ipsec.service remote-fs.target nss-lookup.target
RequiresMountsFor=/data

[Service]
TimeoutStopSec=120
