# etc/systemd/system/memcached.service.d/memcached.conf
[Unit]
After=syslog.target network.target
After=unbound.service
After=network-online.target
After=nss-lookup.target
After=nfs-client.target
Wants=network-online.target
